# 50-elphel-static-nodes.rules
ACTION!="add", GOTO="static_nodes_end"

ACTION=="add", RUN+="/usr/bin/rsync -a /lib/udev/devices/ /dev/"

LABEL="static_nodes_end"
